# hotreddit.timer
[Unit]
Description=Timer for hotreddit.

[Timer]
OnUnitActiveSec=2h
Unit=hotreddit.service

[Install]
WantedBy=timers.target
