# hotreddit.timer
[Unit]
Description=Timer for hotreddit.

[Timer]
OnBootSec=30sec
OnUnitActiveSec=15min
Unit=hotreddit.service

[Install]
WantedBy=timers.target
